# deploy/non-docker/lumina-frontend.service
[Unit]
Description=Lumina Next.js frontend
After=network-online.target lumina-api.service
Wants=network-online.target

[Service]
Type=simple
User=lumina
Group=lumina
WorkingDirectory=/opt/lumina/repo/frontend
EnvironmentFile=/opt/lumina/shared/.env
Environment=PATH=/opt/lumina/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env npm run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
