# deploy/non-docker/lumina-worker.service
[Unit]
Description=Lumina background worker
After=network-online.target docker.service lumina-api.service
Wants=network-online.target

[Service]
Type=simple
User=lumina
Group=lumina
WorkingDirectory=/opt/lumina/repo
EnvironmentFile=/opt/lumina/shared/.env
Environment=PATH=/opt/lumina/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env uv run --env-file /opt/lumina/shared/.env python3 scripts/worker_with_timeout.py --import-modules commands --max-tasks ${SURREAL_COMMANDS_MAX_TASKS}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
