# deploy/non-docker/surrealdb.service
[Unit]
Description=SurrealDB for Lumina
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lumina
Group=lumina
EnvironmentFile=/opt/lumina/.env
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/surreal start --log info --user ${SURREAL_USER} --pass ${SURREAL_PASSWORD} --bind 127.0.0.1:8000 rocksdb:/var/lib/lumina/surrealdb/mydatabase.db
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
